# mine.service
[Unit]
Description=mine

[Service]
ExecStart=/home/$USER/mine/vrm.sh

[Install]
WantedBy=multi-user.target
