# trtl.service
[Unit]
Description=mine

[Service]
ExecStart=/home/miner/mine/./violetminer

[Install]
WantedBy=multi-user.target
