# trtl.service
[Unit]
Description=mine

[Service]
ExecStart=/home/mine/trtl.sh

[Install]
WantedBy=multi-user.target
