# vrm.service
[Unit]
Description=mine

[Service]
ExecStart=/home/mine/vrm.sh

[Install]
WantedBy=multi-user.target
