# vrm.service
[Unit]
Description=mine

[Service]
ExecStart=/home/miner/mine/./vrm.sh

[Install]
WantedBy=multi-user.target
