# xmg.service
[Unit]
Description=mine

[Service]
ExecStart=/home/miner/mine/./xmg.sh

[Install]
WantedBy=multi-user.target
