# xmg.service
[Unit]
Description=mine

[Service]
ExecStart=/home/mine/xmg.sh

[Install]
WantedBy=multi-user.target
